# templates/dehydrated-weekly.service
[Unit]
Description=Updates certificates with Dehydrated

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dehydrated-weekly
